# deploy/morseweb-backup.timer
[Unit]
Description=Nightly morseWeb snapshot to S3

[Timer]
OnCalendar=*-*-* 09:10:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
